# infra/argocd/overlays/minikube/io.nexus.kafka-port-forward.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent — persistent kubectl port-forward for Kafka external listener.

  Forwards two ports to the Kafka NodePort services on 127.0.0.1 (localhost):

    127.0.0.1:32092 → nexus-kafka-external-bootstrap:9093  (bootstrap)
    127.0.0.1:32093 → nexus-combined-external-0:9093       (broker 0)

  The Kafka external listener advertises 127.0.0.1 as the broker address, so
  both the bootstrap and broker connections resolve through these port-forwards.
  This makes KAFKA_BROKERS=localhost:32092 work end-to-end from the macOS host
  without any additional configuration (no sudo, no loopback alias required).

  Install (once):
    cp infra/argocd/overlays/minikube/nexus-kafka-forward.sh \
       ~/.local/bin/nexus-kafka-forward.sh
    chmod +x ~/.local/bin/nexus-kafka-forward.sh
    cp infra/argocd/overlays/minikube/io.nexus.kafka-port-forward.plist \
       ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/io.nexus.kafka-port-forward.plist

  Verify:
    launchctl list io.nexus.kafka-port-forward
    nc -z 192.168.49.2 32092 && echo reachable
    tail -f /tmp/kafka-port-forward.log

  Uninstall:
    launchctl unload -w ~/Library/LaunchAgents/io.nexus.kafka-port-forward.plist
    rm ~/Library/LaunchAgents/io.nexus.kafka-port-forward.plist
    rm ~/.local/bin/nexus-kafka-forward.sh
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.nexus.kafka-port-forward</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/agent/.local/bin/nexus-kafka-forward.sh</string>
  </array>

  <!-- Restart automatically whenever the script exits. -->
  <key>KeepAlive</key>
  <true/>

  <!-- Start on login. -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Back off 10 s between restarts. -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/kafka-port-forward.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/kafka-port-forward.log</string>
</dict>
</plist>
