# infra/argocd/overlays/minikube/io.nexus.minikube-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon — keeps `minikube tunnel` running as root permanently.

  minikube tunnel adds a host route for the minikube network (192.168.49.0/24)
  so that 192.168.49.2 (the minikube node IP) is reachable directly from the
  macOS host — no port-forward required.

  Install (once):
    sudo cp infra/argocd/overlays/minikube/io.nexus.minikube-tunnel.plist \
           /Library/LaunchDaemons/
    sudo launchctl load -w /Library/LaunchDaemons/io.nexus.minikube-tunnel.plist

  Verify:
    sudo launchctl list io.nexus.minikube-tunnel
    tail -f /tmp/minikube-tunnel.log

  Uninstall:
    sudo launchctl unload -w /Library/LaunchDaemons/io.nexus.minikube-tunnel.plist
    sudo rm /Library/LaunchDaemons/io.nexus.minikube-tunnel.plist

  Once the daemon is running the following addresses become reachable:
    ArgoCD UI  https://192.168.49.2:30443
    Kafka CLI  192.168.49.2:32092  (NodePort bootstrap)
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.nexus.minikube-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/minikube</string>
    <string>tunnel</string>
    <string>--cleanup</string>
  </array>

  <!-- Restart immediately if the process exits. -->
  <key>KeepAlive</key>
  <true/>

  <!-- Start when the daemon is loaded (boot / launchctl load). -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Back off 10 s between restarts to avoid tight crash loops. -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/minikube-tunnel.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/minikube-tunnel.log</string>
</dict>
</plist>
